# target/linux/ar71xx/base-files/etc/defconfig/linino-tian/init/mcuio
#!/bin/sh /etc/rc.common
# Copyright (C) 2008 OpenWrt.org

START=11
STOP=11

start() {
    ldattach 29 /dev/ttyACM0
}

stop() {
    killall ldattach
}
