# package/linino/linino-conf/files/etc/crontabs
# Update every 6 hours
* */6 * * * /bin/opkg update >/dev/null 2>&1
