# package/linino/linino-conf/files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

# wifi-refresh

change-bootloader

wifi-live-or-reset

# Update at boot
sleep 120 && /bin/opkg update >/dev/null 2>&1

exit 0
